# debian/postinst
#!/bin/bash
systemctl enable /opt/tanc/systemd/serial-comm.service
