# debian/postinst
#!/bin/bash
systemctl enable /opt/tanc/systemd/tanc-serial-comm.service
systemctl restart tanc-serial-comm
